# supervisor.d/fail2ban.conf
[program:fail2ban]
command = /usr/local/bin/fail2ban-supervisor.sh
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
