# supervisor.d/rsyslog.conf
[program:rsyslogd]
command = rsyslogd -n
